# system/services/mamao-ups.service
[Unit]
Description=MaMaoHub UPS Power Monitor
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/mamao-ups-monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
